# test.service
[Unit]
Description=script for test
Wants=test.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test.sh